# custom/supervisord.d/obsidian-sync.conf
; =============================================================================
; obsidian-sync — official Obsidian Sync CLI (`ob sync --continuous`)
; =============================================================================

[program:obsidian-sync]
command=/opt/hermes-suite/bin/obsidian-sync.sh
directory=/workspace
user=hermes
autostart=true
autorestart=true
priority=60
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/opt/data",HERMES_HOME="/opt/data",OBSIDIAN_VAULT_PATH="/workspace/obsidian",OBSIDIAN_SYNC_ENABLED="true",PATH="/usr/local/bin:/opt/hermes/.venv/bin:/usr/local/sbin:/usr/sbin:/usr/bin:/sbin:/bin"
